# docker-image/dev-toolchain/services.d/dockerd/run
#!/command/execlineb -P
#

s6-envdir /var/run/s6/container_environment
/usr/local/bin/dockerd --host=unix:///var/run/docker.sock --iptables=false
