# docker-image/appbox/services.d/dockerd/run
#!/command/execlineb -P
#

s6-envdir /var/run/s6/container_environment
/usr/bin/dockerd
