# docker-image/dev-toolchain/services.d/sshd/finish
#!/command/execlineb -S0
#

/bin/echo "[sshd] Service stopped"
