# docker-image/dev-toolchain/services.d/vscode/finish
#!/bin/sh
#

echo "[vscode] Service stopped"
