# docker-image/dev-toolchain/services.d/vscode/run
#!/bin/sh
#
# VSCode tunnel service - auto-starts when code CLI is installed

# Wait for code CLI to be installed (check every 10 seconds)
while [ ! -f /usr/local/bin/code ]; do
  echo "[vscode] VSCode CLI not installed yet, waiting..."
  echo "[vscode] Install it with: /usr/local/lib/dev-tools/install-scripts/install-vscode.sh"
  sleep 10
done

echo "[vscode] VSCode CLI detected, starting tunnel..."

# Set tunnel name from environment variable or use hostname as fallback  
TUNNEL_NAME="${VSCODE_TUNNEL_NAME:-$(hostname)}"

# Additional tunnel options
TUNNEL_ARGS="--accept-server-license-terms"

# Add name if specified
if [ -n "$TUNNEL_NAME" ]; then
  TUNNEL_ARGS="$TUNNEL_ARGS --name \"$TUNNEL_NAME\""
fi

# Add verbose logging if enabled
if [ "${VSCODE_TUNNEL_VERBOSE:-false}" = "true" ]; then
  TUNNEL_ARGS="$TUNNEL_ARGS --verbose"
fi

echo "[vscode] Starting tunnel with name: $TUNNEL_NAME"
echo "[vscode] Tunnel args: $TUNNEL_ARGS"
exec /usr/local/bin/code tunnel $TUNNEL_ARGS
