# docker-image/dev-toolchain/services.d/vscode/run
#!/command/execlineb -P
#

s6-envdir /var/run/s6/container_environment
/usr/local/bin/code tunnel --accept-server-license-terms
