# docker-image/dev-toolchain/services.d/vscode/run
#!/bin/sh
#
# VSCode tunnel service - auto-starts when code CLI is installed

# Wait for code CLI to be installed (check every 10 seconds)
while [ ! -f /usr/local/bin/code ]; do
  echo "[vscode] VSCode CLI not installed yet, waiting..."
  echo "[vscode] Install it with: /usr/local/lib/dev-tools/install-scripts/install-vscode.sh"
  sleep 10
done

echo "[vscode] VSCode CLI detected, starting tunnel..."

# Set tunnel name from environment variable or use hostname as fallback
TUNNEL_NAME="${VSCODE_TUNNEL_NAME:-$(hostname)}"

echo "[vscode] Starting tunnel with name: $TUNNEL_NAME"
echo "[vscode] Environment variables:"
echo "[vscode]   VSCODE_TUNNEL_NAME: ${VSCODE_TUNNEL_NAME:-<not set>}"
echo "[vscode]   VSCODE_TUNNEL_VERBOSE: ${VSCODE_TUNNEL_VERBOSE:-<not set>}"

# Build command with proper argument handling
if [ "${VSCODE_TUNNEL_VERBOSE:-false}" = "true" ]; then
  exec /usr/local/bin/code tunnel --accept-server-license-terms --name "$TUNNEL_NAME" --verbose
else
  exec /usr/local/bin/code tunnel --accept-server-license-terms --name "$TUNNEL_NAME"
fi
